# dashboard_workers.conf
[program:ifcb_acq]
user=vagrant
umask=002
directory=/home/vagrant/ifcb-dashboard
environment=PYTHONPATH="/home/vagrant/ifcb-dashboard"
command=celery --config=oii.workflow.async_config -A oii.ifcb2.workflow.acq_worker worker -n ifcb_acq
autorestart=true

[program:ifcb_acc]
user=vagrant
umask=002
directory=/home/vagrant/ifcb-dashboard
environment=PYTHONPATH="/home/vagrant/ifcb-dashboard"
command=celery --config=oii.workflow.async_config -A oii.ifcb2.workflow.acc_worker worker -n ifcb_acc -c 4
autorestart=true

[program:ifcb_batch_acc]
user=vagrant
umask=002
directory=/home/vagrant/ifcb-dashboard
environment=PYTHONPATH="/home/vagrant/ifcb-dashboard"
command=celery --config=oii.workflow.async_config -A oii.ifcb2.workflow.batch_acc_worker worker -n ifcb_batch_acc -c 4
autorestart=true
